# shadow_controller/systemd/shadow-controller.service
[Unit]
Description=Shadow Controller — Hermes-powered DJ Bot Operator
After=network-online.target firefox.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__SCRIPT_DIR__
ExecStart=__SCRIPT_DIR__/venv/bin/python -m shadow_controller.main
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=120

# Environment
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shadow-controller

[Install]
WantedBy=multi-user.target